# docs/examples/com.kuroyk.ibkr-daily-runner.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kuroyk.ibkr-daily-runner</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>PLACEHOLDER_PROJECT_ROOT/scripts/ibkr_local_daily_runner.sh</string>
    <string>--log-root=logs/ibkr_daily</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>16</integer>
    <key>Minute</key>
    <integer>10</integer>
  </dict>
  <key>WorkingDirectory</key>
  <string>PLACEHOLDER_PROJECT_ROOT</string>
  <key>Disabled</key>
  <true/>
</dict>
</plist>
